# services/transmission/run
#!/bin/bash

exec /sbin/setuser abc /usr/bin/transmission-daemon -f -g /config
